# rootfs/etc/services.d/retroassembly/run
#!/usr/bin/with-contenv bashio
# Setup data links
mkdir -p /data/retroassembly
if [ -d "/app/data" ]; then rm -rf /app/data; fi
ln -s /data/retroassembly /app/data

cd /app
exec node dist/server/serve.js
